# setup/www/papertrail_remote_syslog.conf
[program:papertrail_remote_syslog]
command={{ command }}
user = {{ user }}
group = {{ group }}
autostart=true
stdout_logfile = {{ logfile }}
redirect_stderr=true
stopasgroup=true
killasgroup=true
